# config/service/php56/run
#!/bin/bash

exec /usr/sbin/php5-fpm -F > /dev/null 2>&1